# deploy/gunicorn-laman.socket
[Unit]
Description=Gunicorn socket for LAMAN

[Socket]
ListenStream=/run/gunicorn/laman.sock

[Install]
WantedBy=sockets.target
